# src/incBackup.service
[Unit]
Description=Incremental Backup

[Service]
ExecStart=backup

[Install]
WantedBy=multi-user.target 
